# launchd/com.darb-solar.intraday.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.darb-solar.intraday</string>
  <key>WorkingDirectory</key>
  <string>PROJECT_ROOT</string>
  <key>ProgramArguments</key>
  <array>
    <string>PROJECT_ROOT/.venv/bin/python</string>
    <string>PROJECT_ROOT/scripts/sync_fusionsolar_intraday.py</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>StandardOutPath</key>
  <string>PROJECT_ROOT/data/logs/intraday.log</string>
  <key>StandardErrorPath</key>
  <string>PROJECT_ROOT/data/logs/intraday.err.log</string>
</dict>
</plist>
